# debian/postinst
#!/bin/bash
# postinst script for almidon
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package



if [ "$1" = "configure" ];then #Inicio del 1º if

             # Web dir 
             if [ -d /var/www ]; then #Inicio del 2º if
              CONTENTDIR=/var/www
             else
                if [ -d /srv/www ]; then #Inicio del 3º if
                 CONTENTDIR=/srv/www
                fi #Fin del 3º if
             fi #Fin del 2º if

             ALMIDONDIR=$CONTENTDIR/almidon
             if [ ! -d $ALMIDONDIR ];then #Inicio del 4º if
              mkdir $ALMIDONDIR
             fi #Fin del 4º if
                     
            if [ -d /etc/almidon ]; then #Inicio del 5º if
             if [ ! -x /etc/almidon/config.sh ];then #Inicio del 6º if
              chmod +x /etc/almidon/config.sh
             fi #Fin del 6º if
 
             if [ ! -x /etc/almidon/demo-setup.sh ];then #Inicio del 7º if
              chmod +x /etc/almidon/demo-setup.sh
             fi #Fin del 7º if

             if [ ! -x /etc/almidon/site-setup.sh ];then #Inicio del 8º if
              chmod +x /etc/almidon/site-setup.sh
             fi #Fin del 8º if

             if [ ! -x /etc/almidon/phpdoc-almidon.sh ];then #Inicio del 9º if
              chmod +x /etc/almidon/phpdoc-almidon.sh
             fi #Fin del 9º if
 

            fi #Fin del 5º if
                                
            if [ -f /sbin/service ]; then #Inicio del 10º if
	      /sbin/service apache2 condrestart > /dev/null 2>&1 || :
	    else
                 if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then #Inicio del 11º if	  
                  invoke-rc.d apache2 restart
                 else 
                     if [ -f /etc/init.d/apache2 ];then #Inicio del 12º if
                      /etc/init.d/apache2 restart
                     fi #Fin del 12º if
                fi  #Fin del 11º if
            fi  #Fin del 10º if                 
fi #Fin del 1º if

#if [ "$1" = "abort-upgrade" ];then
#fi

#if [ "$1" = "abort-upgrade" ];then
#fi

#if [ "$1" = "abort-remove" ];then
#fi

#if [ "$1" = "abort-deconfigure" ];then
#fi
    



# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


